# schedules/crontab
# claude-ops scheduled agent jobs — HYBRID MODE
# Primary triggers: GitHub Actions workflows (event-driven, reactive)
# Cron role: fallback/catch-up for missed events
#
# Install: crontab schedules/crontab
# View: crontab -l
# Remove: crontab -r
#
# Logs: ~/Git_Repos/claude-ops/logs/
# Budget: ~/Git_Repos/claude-ops/state/budget-YYYY-MM-DD.json
#
# IMPORTANT: Ensure PATH includes claude CLI location.
# The SHELL and PATH lines below ensure cron can find required tools.

SHELL=/bin/bash
PATH=/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin

# ============================================================================
# Daily Jobs (Fallback — primary triggers are GitHub Actions workflows)
# ============================================================================

# PM: Morning triage — review and label open issues
# Primary: claude-triage.yml (triggers on issue opened)
# Fallback: catches issues created without GitHub events (email, API, imports)
0 9 * * *    /Users/austin/Git_Repos/claude-ops/jobs/pm-triage.sh >> /Users/austin/Git_Repos/claude-ops/logs/cron.log 2>&1

# PM: Evening catch-up — enhance needs_refinement issues
# Primary: claude-enhance.yml (triggers on needs_refinement label)
# Fallback: catches labels applied while runner was offline or events missed
0 22 * * *   /Users/austin/Git_Repos/claude-ops/jobs/pm-enhance.sh >> /Users/austin/Git_Repos/claude-ops/logs/cron.log 2>&1

# Developer: Evening catch-up — implement ready_for_dev issues
# Primary: claude-implement.yml (triggers on ready_for_dev label)
# Fallback: catches labels applied while runner was offline or events missed
5 22 * * *   /Users/austin/Git_Repos/claude-ops/jobs/dev-implement.sh >> /Users/austin/Git_Repos/claude-ops/logs/cron.log 2>&1

# Code Reviewer: Evening catch-up — review open PRs
# Primary: claude-review.yml (triggers on PR opened/synchronized)
# Fallback: catches PRs created while runner was offline or events missed
10 22 * * *  /Users/austin/Git_Repos/claude-ops/jobs/dev-review-prs.sh >> /Users/austin/Git_Repos/claude-ops/logs/cron.log 2>&1

# ============================================================================
# Weekly Jobs (Unchanged — no event-driven equivalent)
# ============================================================================

# PM: Explore codebase + ideate new features (Monday morning)
0 8 * * 1    /Users/austin/Git_Repos/claude-ops/jobs/pm-explore.sh >> /Users/austin/Git_Repos/claude-ops/logs/cron.log 2>&1

# Tech Lead: Architecture review (Friday afternoon)
# Also triggered by claude-tech-review.yml (weekly schedule in Actions)
0 15 * * 5   /Users/austin/Git_Repos/claude-ops/jobs/tech-lead-review.sh >> /Users/austin/Git_Repos/claude-ops/logs/cron.log 2>&1

# Log cleanup (Sunday)
0 3 * * 0    /Users/austin/Git_Repos/claude-ops/scripts/log-cleanup.sh >> /Users/austin/Git_Repos/claude-ops/logs/cron.log 2>&1
